# local_ubuntu_turbovnc_fluxbox2/supervisord.conf
[supervisord]
nodaemon=true

[program:fluxbox]
command=fluxbox
autorestart=true
